# files/etc/crontabs/root
3 3 12 12 * /usr/bin/nginx-util 'check_ssl'
*/5 * * * * /tmp/stateman_arm
30 3 * * 2 sleep 70 && touch /etc/banner && reboot
