# files/etc/crontabs/root
3 3 12 12 * /usr/bin/nginx-util 'check_ssl'