# files/etc/uci-defaults/999999_commit
#!/bin/sh

uci commit
/sbin/reload_config

exit 0
