# files/etc/uci-defaults/99_0_preinit
#!/bin/sh

/etc/init.d/cron enable
/etc/init.d/cron start
