# files/etc/uci-defaults/99_3_setup_wan
#!/bin/sh

uci set network.wan.metric=10
uci set network.wan.device='eth1'
uci set network.wan6.device='eth1'

exit 0
