# files/etc/rc.local
printf "AT+QPOWD=1\r" > /dev/ttyUSB2
irqbalance --oneshot
sleep 60
/usr/bin/run_stateman.sh
